# infra/systemd/ea-api.service
[Unit]
Description=Waywarden API
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/waywarden
EnvironmentFile=%h/waywarden/.env
ExecStart=%h/.local/bin/uv run uvicorn waywarden.app:app --app-dir src --host 0.0.0.0 --port 8000
Restart=on-failure

[Install]
WantedBy=default.target
